# packages/agent/agent.service
[Unit]
Description=MCPanel Agent - Minecraft Server Management Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/Minecraft-Server-Management/packages/agent
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcpanel-agent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/minecraft-server
PrivateTmp=true

# Environment
EnvironmentFile=/root/Minecraft-Server-Management/packages/agent/.env

[Install]
WantedBy=multi-user.target
